# supervisor/conf.d/syslog-ng.conf
[program:syslog-ng]
command=syslog-ng 
    -F
autostart=true
autorestart=true
 
